# rpm/qmlpinquery.service
[Unit]
Description=Pin query app
After=xorg.target
Requires=dbus.socket xorg.target 

[Service]
Type=simple
ExecStart=/usr/bin/qmlpinquery

